# code-atlas-search-api.service.example
[Unit]
Description=Code Atlas — search API and web dashboard (Qdrant RAG)
After=network.target

[Service]
Type=simple
User=YOUR_LINUX_USER
Group=YOUR_LINUX_GROUP
WorkingDirectory=/path/to/code-atlas

# Load secrets from a root-owned file: chmod 600 /etc/code-atlas.env
EnvironmentFile=-/etc/code-atlas.env

Environment=PYTHONPATH=/path/to/code-atlas
Environment=QDRANT_PATH=/path/to/code-atlas/data/qdrant_db

# Bind to localhost; put nginx/caddy in front for TLS and access control.
ExecStart=/path/to/code-atlas/.venv/bin/python3 scripts/start_api.py --host 127.0.0.1 --port 8888

Restart=on-failure
RestartSec=5

# Hardening (adjust if you hit limits on large repos)
LimitNOFILE=65535

[Install]
WantedBy=multi-user.target
